# root/etc/services.d/angular/run
#!/usr/bin/with-contenv bash

cd /app || exit

if [ ! -f "/app/angular.json" ]; then
	npm install
	ng new test
	mv ./test/* ./
	rm -rf ./test
fi

exec \
	s6-setuidgid abc /usr/bin/ng serve
